# powerblock-kickstart
#!/system/bin/sh

### BEGIN INIT INFO
# Exec: boot
# Type: permanent
### END INIT INFO

nohup /system/bin/powerblockservice &>/dev/null &
